# broker/contrib/tunneldigger.upstart
description "tunneldigger tunnelling network daemon using l2tpv3"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec /home/tunneldigger/env/bin/python /home/tunneldigger/tunneldigger/broker/l2tp_broker.py /home/tunneldigger/tunneldigger/broker/l2tp_broker.cfg
