# broker/contrib/tunneldigger.upstart
description "tunneldigger tunnelling network daemon using l2tpv3"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

chdir /home/tunneldigger/tunneldigger
exec /home/tunneldigger/env/bin/python -m broker.main /home/tunneldigger/tunneldigger/broker/l2tp_broker.cfg
